# hostvm/cron.d/nutch
# run container once a week
0 0 * * 0 root echo "docker run -it -e ELASTICSEARCH_ENDPOINT=$ELASTICSEARCH_ENDPOINT -e ELASTICSEARCH_HOST=$ELASTICSEARCH_HOST -e ELASTICSEARCH_PORT=$ELASTICSEARCH_PORT --rm jhulibraries/sheridan-libraries-nutch:initialpush /opt/nutch/conf/command.sh " >> /var/log/cron.log 2>&1
# Don't remove the empty line at the end of this file. It is required to run the cron job
