# nutch1/conf/profile.d/nutch.sh
export JAVA_HOME=/usr/lib/jvm/java
export NUTCH_HOME=/opt/nutch
export NUTCH_JAVA_HOME="${JAVA_HOME}"
export NUTCH_RUNTIME="${NUTCH_HOME}/runtime/local"
export NUTCH_LOG_DIR="${NUTCH_RUNTIME}/logs"
export NUTCH_CONF_DIR="${NUTCH_RUNTIME}/conf"
export PATH="${PATH}:${NUTCH_RUNTIME}/bin"
